# jupyterhub/etc/profile.d/python.sh
source scl_source enable rh-python36

if [ -f /opt/app-root/bin/activate ]; then
    source /opt/app-root/bin/activate
fi
